# deb/Factor/zabbix-agent2/DEBIAN/prerm
#!/bin/bash

systemctl stop zabbix_agent2
rm -R /etc/zabbix_agent2
systemctl disable zabbix_agent2
systemctl daemon-reload
systemctl reset-failed
